# infra/systemd/aistareco-web-star.service.example
# /etc/systemd/system/aistareco-web-star.service
# apps/web-star — Next 16 standalone（根路径，默认端口 3014）

[Unit]
Description=AI Star Eco Web Star
After=network-online.target aistareco-server.service
Wants=network-online.target

[Service]
Type=simple
User=ecs-user
WorkingDirectory=/opt/ai-star-eco/web-star/apps/web-star
EnvironmentFile=/etc/aistareco/web-star.env
ExecStart=/usr/bin/node /opt/ai-star-eco/web-star/apps/web-star/server.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
