# infra/systemd/aistareco-web-music.service.example
# /etc/systemd/system/aistareco-web-music.service
# apps/web-music — Next 16 standalone（根路径）

[Unit]
Description=AI Star Eco — Next.js Web-Music (root path)
After=network.target aistareco-server.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ai-star-eco/repo/apps/web-music/.next/standalone

Environment=NODE_ENV=production
Environment=PORT=3010
Environment=HOSTNAME=127.0.0.1

ExecStart=/root/.nvm/versions/node/v24.14.1/bin/node apps/web-music/server.js

Restart=on-failure
RestartSec=5s
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
